# nullifier-ingest/docs/nullifier-query-server.service
[Unit]
Description=Nullifier PIR server (nf-server)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/zally-pir/nullifier-ingest
ExecStart=/opt/nullifier-ingest/nf-server serve \
    --pir-data-dir /root/zally-pir/nullifier-ingest/pir-data \
    --data-dir /root/zally-pir/nullifier-ingest \
    --port 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
