# docs/nullifier-query-server.service
[Unit]
Description=Nullifier PIR server (nf-server)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nf-ingest
ExecStart=/opt/nf-ingest/nf-server serve \
    --pir-data-dir /opt/nf-ingest/pir-data \
    --data-dir /opt/nf-ingest \
    --port 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
