# docs/svoted.service
[Unit]
Description=Shielded-Vote SDK chain (svoted)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/shielded-vote
EnvironmentFile=-/etc/default/svoted
ExecStart=/opt/shielded-vote/svoted start --home /opt/shielded-vote/.svoted
Restart=on-failure

[Install]
WantedBy=multi-user.target
