# docs/svoted-val1.service
[Unit]
Description=Shielded-Vote SDK chain - Validator 1 (svoted-val1)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/shielded-vote
Environment="SVOTE_PIR_URL=http://localhost:3000"
ExecStart=/opt/shielded-vote/svoted start \
    --home /opt/shielded-vote/.svoted-val1 \
    --serve-ui --ui-dist /opt/shielded-vote/ui/dist
Restart=on-failure

[Install]
WantedBy=multi-user.target
