# docs/svoted-archive.service
[Unit]
Description=Shielded-Vote SDK chain - Archive node (svoted-archive)
After=svoted-val1.service
Requires=svoted-val1.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/shielded-vote
ExecStart=/opt/shielded-vote/svoted start --home /opt/shielded-vote/.svoted-archive
Restart=on-failure

[Install]
WantedBy=multi-user.target
