# docs/svoted-val3.service
[Unit]
Description=Shielded-Vote SDK chain - Validator 3 (svoted-val3)
After=network.target svoted-val1.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/shielded-vote
Environment="SVOTE_PIR_URL=http://localhost:3000"
ExecStart=/opt/shielded-vote/svoted start --home /opt/shielded-vote/.svoted-val3
Restart=on-failure

[Install]
WantedBy=multi-user.target
