# docs/zallyd.service
[Unit]
Description=Zally SDK chain (zallyd)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/zally-chain
ExecStart=/opt/zally-chain/zallyd start --home /opt/zally-chain/.zallyd
Restart=on-failure

[Install]
WantedBy=multi-user.target
